# units/borg-backup.timer
[Unit]
Description=Daily Backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=1hour
Persistent=true

[Install]
WantedBy=timers.target
