# units/opt-extra_storage.mount
[Unit]
Description=Extra storage

[Mount]
What=/dev/mapper/extra-storage
Where=/opt/extra_storage
Type=btrfs
