# units/opt-extra_storage.automount
[Unit]
Description=Extra storage

[Automount]
Where=/opt/extra_storage
TimeoutIdleSec=10mins

[Install]
WantedBy=multi-user.target
